# rc.local
#!/bin/sh

iptables -F

iptables -A FORWARD -i wg0 -j ACCEPT

iptables -A FORWARD -o wg0 -j ACCEPT

iptables -t nat -A POSTROUTING -o eth0 -j FULLCONENAT

iptables -t nat -A PREROUTING -i eth0 -j FULLCONENAT

iptables -L

exit
